# cron.daily
#!/bin/sh

python3 /opt/whitelist/scripts/whitelist.py
